# rtorrent-flood/root/etc/services.d/flood/run
#!/bin/sh

cd /app/flood
exec npm start
